# crontab
# /etc/crontab

SHELL=/bin/sh

* * * * * root /home/app_data/docker/openvpn-monitor/scripts/server_status.sh
